# disk-standby.service
[Unit]
Description=Disk Standby Timeout config load

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/sbin/disk-standby --load
